# systemd/hs-fan-mode-init.service
[Unit]
Description=HelloSkyy Supermicro H12 Fan Mode Initialization
After=multi-user.target
Wants=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/hs-fan-mode-init.sh
# Optional: override defaults without editing the script:
# Environment=FAN_MODE=0x04
# Environment=IPMITOOL_BIN=/usr/bin/ipmitool
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

